# haspd/DEBIAN/postinst
#!/bin/sh

systemctl enable aksusbd hasplm
systemctl restart aksusbd hasplm

exit 0
