# haspd/DEBIAN/preinst
#!/bin/sh

case $1 in
  upgrade) systemctl stop aksusbd hasplm
           systemctl disable aksusbd hasplm
           ;;
esac

exit 0
