# haspd/DEBIAN/prerm
#!/bin/bash

systemctl stop aksusbd hasplm
systemctl disable aksusbd hasplm
