# haspd/DEBIAN/prerm
#!/bin/bash

	/etc/init.d/haspd condstop || : 
	test -x /usr/sbin/update-rc.d && update-rc.d -f haspd remove || :